# server/dali-api.service
[Unit]
Description=Dali Party Banner API (FastAPI + rembg)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/dalipart/server
Environment=DALI_DATA_DIR=/var/www/dali-banner-data
Environment=U2NET_HOME=/opt/dali-api/models
Environment=NUMBA_CACHE_DIR=/opt/dali-api/cache
ExecStart=/opt/dali-api/venv/bin/uvicorn app:app --host 127.0.0.1 --port 8019 --workers 1
Restart=on-failure
RestartSec=3

[Install]
WantedBy=multi-user.target
